# unattended-upgrades.timer
[Unit]
Description=Daily (03:00 AM) unattended upgrades

[Timer]
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=timers.target
